# luci2-ffwizard/files/etc/uci-defaults/62_ffwizard_kadnode
#!/bin/sh

. /lib/functions.sh

. /etc/ffwizard.d/60-kadnode.sh
