# files/etc/init.d/ffwizard_autoconf
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 Patrick Grimm

START=98
SERVER="ffwizard_autoconf"
BIN="/usr/sbin/ffwizard_autoconf"


start() {
	firstboot=$(uci_get ffwizard autoconf firstboot "0")
	[ "$firstboot" == "0" ] && return
	start-stop-daemon -S -x $BIN \
		-p /var/run/${SERVER}.pid \
		-m -b --

	uci_set ffwizard autoconf firstboot 0
	uci_commit ffwizard
}

restart() {
	start
}
